# report.timer
[Unit]
Description=Timer file for EOI assignment.

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
