# report.service
[Unit]
Description=Service file for EOI assignment.

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/bin/report.py
